# agent/alesof-agent.service
[Unit]
Description=Alesof on-premise monitoring agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=alesof-agent
Group=alesof-agent
EnvironmentFile=/etc/alesof-agent.env
ExecStart=/usr/bin/python3 /opt/alesof-agent/alesof_agent.py
Restart=always
RestartSec=10
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
